# deploy/go2-relay.service
[Unit]
Description=Go2 relay bridge (local server to cloud relay)
After=network-online.target go2-server.service
Wants=network-online.target go2-server.service

[Service]
Type=simple
User=rpi1
WorkingDirectory=/home/rpi1/dog_mcp
EnvironmentFile=/home/rpi1/dog_mcp/.env
ExecStart=/home/rpi1/dog_mcp/venv/bin/python relay_bridge.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
